# deploy/hosted/systemd/codexsun-platform-api.service
[Unit]
Description=CODEXSUN Platform API
After=network-online.target mariadb.service
Wants=network-online.target

[Service]
Type=simple
User=sundar
Group=sundar
WorkingDirectory=/home/sundar/codexsun
EnvironmentFile=/home/sundar/codexsun/.env
ExecStart=/usr/bin/node /home/sundar/codexsun/dist/apps/platform/api/server.js
Restart=always
RestartSec=3
TimeoutStopSec=30
KillSignal=SIGTERM
NoNewPrivileges=true
PrivateTmp=true
StandardOutput=journal
StandardError=journal
SyslogIdentifier=codexsun-platform-api

[Install]
WantedBy=multi-user.target
